# build_dir/target-mips_24kc_musl/root-ar71xx/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/usr/share/dnsmasq/rfc6761.conf
/etc/config/dhcp
/etc/hotplug.d/ntp/25-dnsmasqsec
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
